# deploy/shakshukabrunch.service
[Unit]
Description=Shakshuka Brunch web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/shakshukabrunch
Environment=PORT=3000
Environment=PATH=/home/deploy/.bun/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env bun run start
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
